# etc/supervisor.d/api.conf
[program:api]
command=/api/bin/server
directory=/api
autostart=true
autorestart=true
user=www
